# radiochecker.service
[Unit]
Description=Radio Checker - Monitor Dutch Radio Stations
After=network.target

[Service]
Type=simple
User=schop
WorkingDirectory=/home/schop/radiochecker
ExecStart=/usr/bin/python3 /home/schop/radiochecker/main.py
Restart=always
RestartSec=10
StandardOutput=append:/home/schop/radiochecker/radio.log
StandardError=append:/home/schop/radiochecker/radio.log

# Optional: Limit log file size
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
